# crontab/crontabConfig
0 * * * * zsh /home/master/dotfiles/zsh/functions/crontabHourly.zsh
*/15 * * * * zsh /home/master/dotfiles/zsh/functions/crontabNote.zsh
