# crontab/crontab5am.fish
#!/usr/bin/env fish

export DISPLAY=:0
hi
