# crontab/crontab_config
@hourly bash /home/master/dotfiles/crontab/crontab_hourly
# daily 7pm
0 19 * * * bash /home/master/dotfiles/crontab/crontab_daily
