# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bash/media
source ~/dotfiles/bash/completion

update_music

update_bash_functions

dnf makecache
