# crontab/crontab_config
@hourly bash /home/master/dotfiles/crontab/crontab_hourly
@daily bash /home/master/dotfiles/crontab/crontab_daily
