# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bashrc

update

update_repo

sync_repo

clone_all_users_repo
