# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bashrc

update_music
update
