# crontab/crontab_hourly
#!/bin/bash

source ~/dotfiles/bash/dev

update_all_repo

auto_sync
