# crontab/crontabConfig
0 * * * * zsh /home/master/git/dotfiles/zsh/functions/crontabHourly.zsh
