# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bash/media

update_music

tldr --update

sudo dnf update -y
sudo dnf upgrade -y
dnf makecache

echo $(date) >>~/.local/share/test/crontab_daily.txt
