# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bashrc

update

clone_all_users_repo
