# crontab/crontab_daily
#!/bin/bash

source ~/dotfiles/bash/media

update_music
