# crontab/crontabConfig
@hourly zsh /home/master/dotfiles/zsh/functions/crontabHourly.zsh
@daily zsh /home/master/dotfiles/zsh/functions/crontabDaily.zsh
*/15 * * * * zsh /home/master/dotfiles/zsh/functions/crontab15m.zsh
