# crontab/crontab_weekly
#!/bin/bash

source ~/dotfiles/bashrc

sync_music
