# crontab/crontab_hourly
#!/bin/bash

source ~/dotfiles/bashrc

update_repo
