# crontab/crontab_hourly
#!/bin/bash

source ~/dotfiles/bashrc

update_all_repo

auto_sync
